# build-rcouch-merge/rcouch.conf
[program:rcouch]
command=/couchdb/rel/apache-couchdb/bin/couchdb start
