# rcouch-merged/rcouch.conf
[program:rcouch]
command=/rel/apache-couchdb/bin/couchdb start
