# build-rcouch-stock/rcouch.conf
[program:rcouch]
command=/rcouch/rel/rcouch/bin/rcouch start

